# systemd/system/printer-ui-api.service
[Unit]
Description=Printer UI Update Host API
After=network.target

[Service]
Type=simple
Environment=UPDATE_API_PORT=9876
WorkingDirectory=/opt/printer-ui
ExecStart=/usr/bin/node /opt/printer-ui/host-api/server.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
