# systemd/system/printer-ui-update-api.service
sudo tee /etc/systemd/system/printer-ui-update-api.service > /dev/null <<'EOF'
[Unit]
Description=Printer UI Host Update API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/printer-ui
Environment=UPDATE_API_PORT=9876
ExecStart=/usr/bin/node /opt/printer-ui/host-api/server.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
EOF
