# systemd/system/printer-ui-update-api.service
[Unit]
Description=Printer UI Host Update API
After=network.target

[Service]
Type=simple
User=@APP_USER@
WorkingDirectory=/opt/printer-ui
Environment=UPDATE_API_PORT=9876
ExecStart=/usr/bin/node /opt/printer-ui/host-api/server.js
Restart=always
RestartSec=2
Environment=UPDATE_API_HOST=0.0.0.0

[Install]
WantedBy=multi-user.target
