# systemd/printer-ui-update.service
[Unit]
Description=Printer UI Update (git pull + docker rebuild)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/printer-ui
ExecStart=/opt/printer-ui/scripts/update.sh
TimeoutStartSec=0
StandardOutput=append:/opt/printer-ui/data/update.log
StandardError=append:/opt/printer-ui/data/update.log

[Install]
WantedBy=multi-user.target
