# systemd/printer-ui-electron.service.in
[Unit]
Description=Printer UI Electron Kiosk
After=graphical.target printer-ui.service
Requires=printer-ui.service
Wants=graphical.target

[Service]
Type=simple
User=@APP_USER@
WorkingDirectory=/opt/printer-ui
Environment=DISPLAY=:0
Environment=KIOSK_URL=http://localhost:3000

ExecStart=/opt/printer-ui/electron-app/start.sh

Restart=always
RestartSec=2

[Install]
WantedBy=graphical.target
