# systemd/printer-ui.service.in
[Unit]
Description=Printer UI Backend (Docker Compose)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/printer-ui
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes
TimeoutStartSec=0
User=@APP_USER@

[Install]
WantedBy=multi-user.target
