# systemd/user/printer-ui-electron.service.in
[Unit]
Description=Printer UI Electron Kiosk (User)
After=graphical-session.target printer-ui.service
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/opt/printer-ui
Environment=KIOSK_URL=http://localhost:3000
ExecStart=/opt/printer-ui/electron-app/start.sh
Restart=always
RestartSec=2
TimeoutStartSec=90

[Install]
WantedBy=default.target
