# systemd/user/printer-ui-electron.service.in
[Unit]
Description=Printer UI Electron Kiosk
After=graphical-session.target printer-ui.service
Wants=graphical-session.target

[Service]
Type=simple

# App l채uft IM User-Kontext (wichtig f체r X11/Wayland)
WorkingDirectory=/opt/printer-ui

# Display aus der Desktop-Session
Environment=DISPLAY=:0

# Ziel-URL f체r den Kiosk
Environment=KIOSK_URL=http://localhost:3000

# Electron starten
ExecStart=/opt/printer-ui/electron-app/start.sh

# Stabilit채t
Restart=always
RestartSec=2
TimeoutStartSec=90

[Install]
WantedBy=default.target
